# .config/autostart/diff_checker.desktop
[Desktop Entry]
Name=Diff Checker
Exec=/opt/diff-checker/diff-checker
Terminal=false
Type=Application
X-GNOME-Autostart-enabled=true
