# .config/autostart/mount_archive_share.desktop
[Desktop Entry]
Name=Mount Archive
Exec=gio mount smb://nas-charlie.local/archive
Terminal=false
Type=Application
X-GNOME-Autostart-enabled=true
