# .config/autostart/mount_private_share.desktop
[Desktop Entry]
Name=Mount Private
Exec=gio mount smb://nas-charlie.local/private
Terminal=false
Type=Application
X-GNOME-Autostart-enabled=true
