# .config/autostart/stash.desktop
[Desktop Entry]
Name=Stash
Exec=/home/charlierz/bin/stash-linux
Terminal=false
Type=Application
X-GNOME-Autostart-enabled=true
